# templates/default/redis@.service.erb
[Unit]
Description=Redis (<%= @name %>) persistent key-value database
After=network.target

[Service]
ExecStart=<%= @bin_path %>/redis-server <%= @configdir %>/<%= @name %>.conf --daemonize no
User=<%= @user %>
Group=<%= @group %>
LimitNOFILE=<%= @limit_nofile %>

[Install]
WantedBy=multi-user.target
